# systemd/user/fuse-promised.service
[Unit]
Description=fuse-promise user-session daemon

[Service]
Type=simple
# Placeholder template. The current skeleton serves private status IPC, but the
# FUSE adapter is not implemented yet.
ExecStart=/usr/bin/fuse-promised --foreground
Restart=no

[Install]
WantedBy=default.target
