# systemd/user/fuse-promised.service
[Unit]
Description=fuse-promise user-session daemon

[Service]
Type=simple
# Placeholder template. Do not install until fuse-promised keeps a FUSE mount
# alive; the current skeleton exits with failure in --foreground mode.
ExecStart=/usr/bin/fuse-promised --foreground
Restart=no

[Install]
WantedBy=default.target
